# scripts/com.danielferrari.screentime-ui-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the iPhone Screen Time UI scraper.

  Runs daily at 21:00 (local time). Opens System Settings briefly,
  switches the Device popup to iPhone, scrapes the App & Website
  Activity table, POSTs to /api/screentime/ingest with
  source="mac_ui_iphone". Window steals focus for ~10–15 seconds
  during the scrape.

  Companion to com.danielferrari.screentime-sync.plist (the
  knowledgeC.db Mac collector). Both run; readers in the dashboard
  decide which row wins per (date, app).

  Install:
    1. Edit the placeholders below (REPLACE_*).
    2. Copy:
         cp scripts/com.danielferrari.screentime-ui-sync.plist \
            ~/Library/LaunchAgents/com.danielferrari.screentime-ui-sync.plist
    3. Load it:
         launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.danielferrari.screentime-ui-sync.plist
    4. (Optional) Trigger now:
         launchctl kickstart -k gui/$UID/com.danielferrari.screentime-ui-sync

  Inspect:
    launchctl list | grep screentime-ui
    tail -f /tmp/screentime-ui-sync.log

  Uninstall:
    launchctl bootout gui/$UID ~/Library/LaunchAgents/com.danielferrari.screentime-ui-sync.plist
    rm ~/Library/LaunchAgents/com.danielferrari.screentime-ui-sync.plist

  REQUIRED ONE-TIME SETUP: Accessibility permission for /usr/bin/osascript.
  System Settings → Privacy & Security → Accessibility → "+", press
  Cmd+Shift+G, enter /usr/bin/osascript, add it, toggle on.
  TCC tracks the leaf binary that opens the AX connection — granting
  Accessibility to the wrapper (env, npx, tsx, node) does not propagate.

  Also: enable "Share Across Devices" in Settings → Screen Time on both
  the Mac and the iPhone, otherwise the device popup will not have an
  iPhone entry.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.danielferrari.screentime-ui-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>npx</string>
    <string>--yes</string>
    <string>tsx</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/Personal Dashboard/scripts/screentime-ui-sync.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH/Personal Dashboard</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>SCREENTIME_INGEST_URL</key>
    <string>REPLACE_WITH_HTTPS_URL/api/screentime/ingest</string>
    <key>SCREENTIME_INGEST_SECRET</key>
    <string>REPLACE_WITH_SHARED_SECRET</string>
    <key>SCREENTIME_TZ</key>
    <string>Australia/Sydney</string>
    <!-- Wait for ~7m58s of HID idle before scraping. -->
    <key>SCREENTIME_UI_IDLE_S</key>
    <string>478</string>
    <!-- Only run outside working hours: before 09:00 or at/after 17:00. -->
    <key>SCREENTIME_UI_WORK_START_H</key>
    <string>9</string>
    <key>SCREENTIME_UI_WORK_END_H</key>
    <string>17</string>
  </dict>

  <!-- Every 2 minutes. The scrape itself takes 3-5 minutes (polling
       SwiftUI's lazy-rendering table), but the script gates on idle
       time, working hours, and a 4-hour cooldown — most invocations
       exit in <1 second without doing anything. The scraper uses
       `open -gja` + AX scripting so when it does run, the active Space
       never switches. If the Mac is asleep, launchd fires missed runs
       on next wake.

       Tunables (override in EnvironmentVariables above):
         SCREENTIME_UI_IDLE_S        default 120   (HID-idle seconds)
         SCREENTIME_UI_COOLDOWN_S    default 14400 (4 h between successes)
         SCREENTIME_UI_WORK_START_H  default 9     (skip hours [start, end))
         SCREENTIME_UI_WORK_END_H    default 17
  -->
  <key>StartInterval</key>
  <integer>120</integer>

  <!-- Append (don't overwrite) so the running log accumulates. -->
  <key>StandardOutPath</key>
  <string>/tmp/screentime-ui-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/screentime-ui-sync.log</string>
</dict>
</plist>
